# docker/cron/crontab
# Cryptocurrency rates fetching every 5 minutes
*/5 * * * * cd /var/www/html && php bin/console app:fetch-rates >> /var/log/cron.log 2>&1

# Database cleanup weekly (keep 1 year of data)
0 2 * * 0 cd /var/www/html && php bin/console app:cleanup-rates >> /var/log/cron.log 2>&1

# Cache warmup daily
0 3 * * * cd /var/www/html && php bin/console cache:warmup --env=prod >> /var/log/cron.log 2>&1
